# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:web]
command=uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s --workers %(ENV_WEB_WORKERS)s
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=5

[program:celery-worker]
command=celery -A app.tasks.celery_app worker --loglevel=%(ENV_LOG_LEVEL)s --concurrency=%(ENV_CELERY_CONCURRENCY)s
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=10

[program:celery-beat]
command=celery -A app.tasks.celery_app beat --loglevel=%(ENV_LOG_LEVEL)s
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=10
